# ihub.conf
[program:ihub_flask]
command= python3 /root/ihub/app.py
autostart=true
autorestart=true
stderr_logfile=/root/ihub/app.err.log
stdout_logfile=/root/ihub/app.out.log